# playground/resources/supervisor/programs/seperate_audio_demo.conf
[program:train_word_vectors_demo_ui]
command=opyrator launch-ui app:train_word_vectors --port=7082
directory=%(ENV_RESOURCES_PATH)s/demos/train_word_vectors
redirect_stderr=true
;stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
stdout_logfile=/proc/1/fd/1 ; log error into stdout to see in docker logs
stdout_logfile_maxbytes=0
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startretries=5   ; max # of serial start failures (default 3)

[program:train_word_vectors_demo_api]
command=opyrator launch-api app:train_word_vectors --port=6082
directory=%(ENV_RESOURCES_PATH)s/demos/train_word_vectors
redirect_stderr=true
;stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
stdout_logfile=/proc/1/fd/1 ; log error into stdout to see in docker logs
stdout_logfile_maxbytes=0
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startretries=5   ; max # of serial start failures (default 3)
